# service-slm/scripts/systemd/foundry-slm-canary.service
[Unit]
Description=SLM canary — Tier A throughput/throttle gate + drain liveness (one-shot)
Documentation=file:///srv/foundry/clones/project-totebox/service-slm/scripts/slm-canary.sh
# Catches the regression classes that slipped through for 10 days (2026-06):
#   - throughput collapse (cgroup repack/MemoryMax thrash, 4 -> 0.3 tok/s)
#   - drain stalls (empty-diff OLMo hang; infinite-retry brief)
After=local-slm.service local-doorman.service

[Service]
Type=oneshot
User=mathew
Group=foundry
WorkingDirectory=/srv/foundry
# NOTE (Command): copy the three scripts to a stable deployed path and point
# ExecStart there, or run from the canonical mirror. Reference (archive) path:
ExecStart=/srv/foundry/clones/project-totebox/service-slm/scripts/slm-canary.sh
StandardOutput=journal
StandardError=journal
PrivateTmp=true
# Tunable gates (override here if the hardware baseline changes):
Environment=SLM_MIN_TOKS_PER_SEC=3.5
Environment=SLM_MAX_HIGH_DELTA=500
Environment=SLM_STALE_INFLIGHT_MIN=30

[Install]
# No WantedBy — invoked by the timer only.
