# infrastructure/systemd/mediakit/local-proofreader.service
[Unit]
Description=service-proofreader (operational write-assistant HTTP backend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=foundry
Group=foundry

WorkingDirectory=/opt/mediakit/data/proofreader

Environment="HOME=/opt/mediakit/data/proofreader"
Environment="PROOFREADER_BIND=0.0.0.0:9092"
Environment="LANGUAGETOOL_URL=http://127.0.0.1:8010"
Environment="DOORMAN_URL=http://127.0.0.1:9080"
Environment="PROOFREADER_DOORMAN_MODEL=olmo3"
Environment="PROOFREADER_DOORMAN_TIMEOUT_SECS=360"
Environment="PROOFREADER_DOORMAN_MAX_TOKENS=256"
Environment="RUST_LOG=info"

ExecStart=/opt/mediakit/bin/service-proofreader

Restart=on-failure
RestartSec=5s

MemoryMax=512M
CPUQuota=200%

NoNewPrivileges=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
