# infrastructure/systemd/mediakit/local-fs.service
[Unit]
Description=service-fs Ring 1 WORM Immutable Ledger (per-tenant boundary-ingest)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=foundry
Group=foundry

WorkingDirectory=/opt/mediakit/data/service-fs

Environment="HOME=/opt/mediakit/data/service-fs"
Environment="FS_BIND_ADDR=0.0.0.0:9100"
Environment="FS_MODULE_ID=foundry-workspace"
Environment="FS_LEDGER_ROOT=/opt/mediakit/data/service-fs/ledger"
Environment="RUST_LOG=info"

ExecStart=/opt/mediakit/bin/service-fs

Restart=on-failure
RestartSec=5s

MemoryMax=512M
CPUQuota=200%

NoNewPrivileges=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
