# service-slm/compute/packer/scripts/llama-server.service
[Unit]
Description=llama.cpp OpenAI-compatible inference server (Yo-Yo Tier B)
After=vllm-weights-prep.service network.target
Requires=vllm-weights-prep.service

[Service]
Type=simple
ExecStart=/opt/llama.cpp/build/bin/llama-server \
    -m /data/weights/olmo-3-32b-think-q3.gguf \
    -ngl 99 \
    --host 127.0.0.1 \
    --port 8000 \
    -c 4096 \
    -a Olmo-3-1125-32B-Think,Olmo-3-32B-Think,olmo-3-32b \
    -np 4 \
    --metrics
Restart=on-failure
RestartSec=15
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
