# service-slm/compute/systemd/nightly-run.timer
[Unit]
Description=Yo-Yo #1 nightly automation (DataGraph rebuild + LoRA training)
Documentation=file:///srv/foundry/clones/project-intelligence/service-slm/scripts/nightly-run.sh

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
