# infrastructure/systemd/orchestration/local-bim-orchestration.service
[Unit]
Description=app-orchestration-bim (PointSav Building Design System Substrate HTTP backend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=foundry
Group=foundry

WorkingDirectory=/opt/mediakit/data/bim-orchestration

Environment="HOME=/opt/mediakit/data/bim-orchestration"
Environment="BIM_BIND=0.0.0.0:9096"
Environment="BIM_VAULT_DIR=/opt/mediakit/data/bim-orchestration/vault"
Environment="BIM_TENANT=woodfine"
Environment="BIM_THEME=pointsav-brand"
Environment="BIM_PUBLIC_URL=https://bim.woodfinegroup.com"
Environment="BIM_DESIGN_SYSTEM_DIR=/opt/mediakit/data/bim-orchestration/design-system"
Environment="BIM_STATIC_DIR=/opt/mediakit/data/bim-orchestration/static"
Environment="FS_ENDPOINT=http://127.0.0.1:9100"
Environment="FS_MODULE_ID=foundry-workspace"
Environment="DOORMAN_URL=http://127.0.0.1:9080"
Environment="RUST_LOG=info"

ExecStart=/opt/mediakit/bin/app-orchestration-bim

Restart=on-failure
RestartSec=5s

MemoryMax=512M
CPUQuota=200%

NoNewPrivileges=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
