# infrastructure/systemd/local-knowledge-documentation.service
[Unit]
Description=PointSav Documentation Wiki (app-mediakit-knowledge)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=foundry
Group=foundry

WorkingDirectory=/opt/mediakit/data/knowledge

Environment="WIKI_CONTENT_DIR=/opt/mediakit/data/content-wiki-documentation"
Environment="WIKI_BIND=0.0.0.0:9090"
Environment="WIKI_SITE_TITLE=PointSav Documentation Wiki"
Environment="WIKI_CITATIONS_YAML=/opt/mediakit/data/content-wiki-documentation/citations.yaml"
Environment="WIKI_STATE_DIR=/opt/mediakit/data/knowledge/state-documentation"

ExecStart=/opt/mediakit/bin/app-mediakit-knowledge serve

Restart=on-failure
RestartSec=5s

NoNewPrivileges=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
